# ping-chime_1.0-5/DEBIAN/prerm
#!/bin/sh
set -e

# Stop and disable the service
systemctl stop ping-chime
systemctl disable ping-chime

exit 0
